# scripts/etc/supervisor/conf.d/listenLIRC.conf
[program:listenLIRC]
command=sudo /opt/bin/listenLIRC >/dev/null
autostart=true
autorestart=true
redirect_stderr=true
user=pi
stdout_logfile=/dev/null
stderr_logfile=/dev/null